# copy/.config/systemd/user/hourly@.timer
[Unit]
Description=Hourly Timer for %i service

[Timer]
OnCalendar=hourly
Persistent=true
Unit=%i.service

[Install]
WantedBy=default.target
