# copy/.config/systemd/user/monthly@.timer
[Unit]
Description=Monthly Timer for %i service

[Timer]
OnCalendar=monthly
Persistent=true
Unit=%i.service

[Install]
WantedBy=default.target
